# systemd/webfm.service
[Unit]
Description=Web File Manager Pro by LONE FAISAL
Documentation=https://github.com/lonefaisal7/web-file-manager
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/webfm
Environment=PORT=8000
Environment=FLASK_ENV=production
ExecStart=/opt/webfm/.venv/bin/python /opt/webfm/app.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=webfm

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/webfm

[Install]
WantedBy=multi-user.target
